# deploy/supervisor/erptechindia-worker.conf
; Supervisor configuration for the Laravel queue worker.
;
; Install once on the server:
;   sudo cp /var/www/erptechindia.tech/deploy/supervisor/erptechindia-worker.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start erptechindia-worker:*
;
; Check status:
;   sudo supervisorctl status erptechindia-worker:*
;
; After deploys, restart workers so they pick up new code:
;   php artisan queue:restart        (run from the project dir)
;
; This config runs 2 workers in parallel (numprocs=2). They share the
; database queue and process jobs in FIFO order. If any worker crashes,
; supervisor restarts it automatically.

[program:erptechindia-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/erptechindia.tech/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --backoff=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/erptechindia.tech/storage/logs/queue-worker.log
stopwaitsecs=3600
